# service-logger.service
[Unit]
Description=Service Logger
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/service-logger.sh

[Install]
WantedBy=default.target
